# ha-agent.service
# /etc/systemd/system/ha-agent.service
[Unit]
Description=HA Agent Service
After=network.target

[Service]
User=youruser
Group=youruser
WorkingDirectory=/srv/app/ha/current
ExecStart=/srv/app/ha/current/ha-agent
EnvironmentFile=/srv/app/ha/app.env
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target